# update-motd.d/20-date
#!/bin/sh

echo "Can you believe it's already $(date)?!"
